# forgejo/rootfs/etc/s6-overlay/s6-rc.d/postgres/run
#!/command/execlineb -P
# Run postgres as the postgres user with PGDATA=/data/postgres
fdmove -c 2 1
s6-setuidgid postgres
/usr/bin/postgres -D /data/postgres
